# package/third-party/udev/11-mount-usb.rules
ACTION=="add", GOTO="automountusb", KERNEL=="sd[a-z][0-9]", RUN+="/bin/mount -t vfat /dev/%k /mnt/", LABEL="automountusb"
